# bare-metal/stt/config/whisper-idle-monitor.service
[Unit]
Description=Whisper Idle Monitor
After=whisper-api.service

[Service]
Type=simple
ExecStart=%h/.local/bin/whisper-idle-monitor.sh

[Install]
WantedBy=default.target
