# bare-metal/llama-cpp/config/llama-cpp.service
[Unit]
Description=llama.cpp Local LLM Server (Router Mode)
After=network.target

[Service]
Type=simple

# Router mode with model management API
# --models-preset: Load model configurations from presets.ini
# --models-max 1: Only 1 model at a time (4GB VRAM constraint)
# --sleep-idle-seconds 300: Auto-unload after 5 min idle
# Exposes API: POST /models/unload, POST /models/load, GET /models
ExecStart=/bin/bash -c 'if [ -f %h/.zshenv ]; then set -a; source %h/.zshenv; set +a; fi && exec %h/.local/share/ai-stack/llama-cpp/llama.cpp/build/bin/llama-server \
    --models-preset %h/.config/ai-stack/llama-cpp/presets.ini \
    --models-max 1 \
    --sleep-idle-seconds 300 \
    --host 127.0.0.1 \
    --port 7865 \
    --metrics'

Restart=on-failure
RestartSec=5s

# Allow higher priority for inference
Nice=-5

[Install]
WantedBy=default.target