# bare-metal/llama-cpp/config/llama-cpp.service
[Unit]
Description=llama.cpp Local LLM Server
After=network.target

[Service]
Type=simple

# Build with CUDA, LTO, and RTX 3050 architecture (sm_86)
# Uses -hf flag for auto-download from HuggingFace
# Sources ~/.zshenv directly to get HF_TOKEN
ExecStart=/bin/bash -c 'if [ -f %h/.zshenv ]; then set -a; source %h/.zshenv; set +a; fi && exec %h/.local/share/ai-stack/llama-cpp/llama.cpp/build/bin/llama-server \
    -hf unsloth/Qwen3.5-4B-GGUF:Q4_K_M \
    --host 127.0.0.1 \
    --port 7865 \
    --metrics \
    --sleep-idle-seconds 300

Restart=on-failure
RestartSec=5s

# Allow higher priority for inference
Nice=-5

[Install]
WantedBy=default.target